# core/detection_engine/rules/default/registry_persistence.yml
title: Registry Run Key Persistence
id: registry_run_key
description: Detects creation of registry Run/RunOnce keys — common persistence technique.
status: experimental
logsource:
    category: registry_event
    product: windows
detection:
    selection:
        command_line|contains:
            - 'Run\\'
            - 'RunOnce\\'
            - 'CurrentVersion\\Run'
    condition: selection
level: medium
tags:
    - attack.persistence
    - attack.t1547.001
